# crontabs
0 18 * * 0 gsutil -m rsync -r -x '.*\.(jpg|jpeg|png|gif|bmp|tiff|svg|)$' -x '.*/outputs/*' /Users/douglashindson/workspace/foundation gs://allmytext/foundation
0 18 * * 0 gsutil -m rsync -r -x '.*\.(jpg|jpeg|png|gif|bmp|tiff|svg|)$' -x '.*/outputs/*' /Users/douglashindson/workspace/blog gs://allmytext/blog